# vpn-switcher.service
[Unit]
Description=VPN Switcher Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/vpn-switcher
Restart=on-failure
RestartSec=5
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
